# agent/com.selfcontrol.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.selfcontrol.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/self-control-agent</string>
        <string>--service</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/self-control-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/self-control-agent.err</string>
</dict>
</plist>
